# Gig Copilot Bot/launchd/com.neonblonde.gigcopilot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.neonblonde.gigcopilot</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Volumes/VADER/Manifold/Neon_Blonde/Repos/Neon_v2/Gig Copilot Bot/scripts/run_gig_copilot.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/studio_hub/.hermes/logs/gigcopilot.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/studio_hub/.hermes/logs/gigcopilot.log</string>
  <key>WorkingDirectory</key>
  <string>/Volumes/VADER/Manifold/Neon_Blonde/Repos/Neon_v2/Gig Copilot Bot</string>
</dict>
</plist>
